# deploy/tooltally.service
[Unit]
Description=ToolTally Kiosk Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tooltally
WorkingDirectory=/opt/tooltally
ExecStart=/opt/tooltally/.venv/bin/python /opt/tooltally/main.py --config /opt/tooltally/config/settings.yaml
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
